# mailservice.app/templates/default/upstart.mailservice.queue.conf.erb
description "Mailservice Resque workers"

start on starting redis
stop on stopped redis

respawn
respawn limit 10 5
console none

setuid <%= node[:unicorn][:user] %>
setgid <%= node[:unicorn][:group] %>

chdir <%= node[:unicorn][:apps_dir] %>/mailservice.app

exec /usr/local/bin/rake resque:work QUEUE=* PIDFILE=<%= node[:unicorn][:apps_dir] %>/mailservice.app/tmp/pids/resque.pid  RAILS_ENV=production >> <%= node[:unicorn][:apps_dir] %>/mailservice.app/log/resque.log 2>&1
