# vps/templates/default/upstart.app.conf.erb
description "<%= @app_name %> application"

start on (local-filesystems and net-device-up IFACE=eth0)
stop on shutdown

respawn
respawn limit 10 5
console none

setuid <%= node[:unicorn][:user] %>
setgid <%= node[:unicorn][:group] %>
chdir <%= node[:unicorn][:apps_dir] %>/<%= @app_name %>

exec /usr/local/bin/unicorn_rails -c <%= node[:unicorn][:apps_dir] %>/<%= @app_name %>/unicorn.rb -E production 
