# notification-service/templates/default/notification-scheduler-monitor.conf.erb
# upstart-respawn-monitor.conf

description "Notify root about notification-scheduler respawn limit reached"
author "Roman Yepishev <rtg@rtg.in.ua>"

# This is a short-lived job
task

start on stopped notification-scheduler PROCESS=respawn

script
    echo "notification-service reached respawn limit" | mail -s "notification-scheduler crashed" notifications@rammount.com
end script
