# notification-service/templates/default/upstart.notification-loader.conf.erb
description "Notification Loader"

start on starting redis
stop on stopped redis

respawn
respawn limit 10 5
console none

setuid <%= node[:unicorn][:user] %>
setgid <%= node[:unicorn][:group] %>

expect daemon

chdir <%= node[:unicorn][:apps_dir] %>/<%= @app_name %>
env RAILS_ENV=<%= node['notification-service'][:environment] %> 

script
  PIDFILE=<%= node[:unicorn][:apps_dir] %>/<%= @app_name %>/tmp/pids/notification-loader.pid
  echo $$ > $PIDFILE
  exec script/daemon start notification_loader >> <%= node[:unicorn][:apps_dir] %>/<%= @app_name %>/log/notification-loader.log 2>&1
end script

