# notification-service/templates/default/upstart.notification-loader.conf.erb
description "Notification Loader"

start on starting redis
stop on stopped redis

respawn
respawn limit 3 15
console none

setuid <%= node[:unicorn][:user] %>
setgid <%= node[:unicorn][:group] %>

expect daemon

chdir <%= node[:unicorn][:apps_dir] %>/<%= @app_name %>
env RAILS_ENV=<%= node['notification-service'][:environment] %> 

exec /usr/local/bin/ruby script/daemon start notification_loader >> <%= node[:unicorn][:apps_dir] %>/<%= @app_name %>/log/notification-loader.log 2>&1
