# notification-service/templates/default/upstart.notification-service.conf.erb
description "Notification Resque Workers"

start on starting redis
stop on stopped redis

respawn
respawn limit 10 5
console none

setuid <%= node[:unicorn][:user] %>
setgid <%= node[:unicorn][:group] %>

chdir <%= node[:unicorn][:apps_dir] %>/<%= @app_name %>
script
  PIDFILE=<%= node[:unicorn][:apps_dir] %>/<%= @app_name %>/tmp/pids/resque.pid
  echo $$ > $PIDFILE
  exec /usr/local/bin/ruby -S bundle exec rake environment resque:work QUEUE=* RAILS_ENV=<%= node['notification-service'][:environment] %> PIDFILE=$PIDFILE >> <%= node[:unicorn][:apps_dir] %>/<%= @app_name %>/log/resque.log 2>&1
end script
